# 7.4/laravel-nginx-mysql/worker.ini
[program:worker]
command = php /var/www/artisan queue:work
user = www-data