# 7.4/mysql-nginx/master.ini
[supervisord]
logfile=/var/log/supervisord.log 
logfile_maxbytes=50MB        
logfile_backups=10           
loglevel=info                
pidfile=/tmp/supervisord.pid 
nodaemon=true   
user=root            
minfds=1024                 
minprocs=200                

[program:crond]
command=crond -f -L 15
user=root
autostart=true

[program:php]
command = docker-php-entrypoint php-fpm

[program:nginx]
command=nginx -g 'pid /tmp/nginx.pid; daemon off;'
autostart=true
autorestart=true
priority=10
stdout_events_enabled=true
stderr_events_enabled=true

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface


[supervisorctl]
serverurl=unix:///run/supervisor.sock ; use a unix:// URL  for a unix socket

[include]
files = /etc/supervisor.d/*.ini

[unix_http_server]
file=/run/supervisor.sock   ; the path to the socket file
